# .config/systemd/user/dzenify.service
[Unit]
Description=Notification daemon

[Service]
Type=simple
ExecStart=/home/oskar/script/dzenify/dzenify.py

[Install]
WantedBy=multi-user.target
